# scripts/com.iddl.indd-dashboard-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  LaunchAgent for the weekly INDD dashboard data refresh.

  Runs every Monday at 05:00 local time. If the Mac is asleep at that
  time, launchd defers the job until the next wake. If the Mac is off,
  the run is skipped entirely and the next chance is the following Monday.

  Install:
    cp scripts/com.iddl.indd-dashboard-refresh.plist ~/Library/LaunchAgents/
    launchctl load -w ~/Library/LaunchAgents/com.iddl.indd-dashboard-refresh.plist

  One-time TCC grant (required because the project lives in ~/Documents):
    System Settings -> Privacy & Security -> Full Disk Access -> "+"
    -> Cmd+Shift+G -> /bin/bash -> Open -> toggle on.
    Without this, launchd cannot read this script and the run fails
    with "Operation not permitted" (exit 126).

  Inspect:
    launchctl list | grep com.iddl.indd
    tail -F ~/Library/Logs/INDD_dashboard_refresh.log
    tail -F ~/Library/Logs/INDD_dashboard_refresh.{out,err}.log

  Run on demand (without waiting for Monday):
    launchctl start com.iddl.indd-dashboard-refresh

  Disable:
    launchctl unload -w ~/Library/LaunchAgents/com.iddl.indd-dashboard-refresh.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.iddl.indd-dashboard-refresh</string>

    <!--
      Invoked via osascript (do shell script) instead of /bin/bash directly.
      This detaches the spawned shell from launchd's immediate child domain,
      which works around an EDEADLK ("Resource deadlock avoided") that
      mmap() hits when launchd-spawned processes touch files inside an
      iCloud-Drive-synced ~/Documents (FileProvider extension contention).
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/osascript</string>
        <string>-e</string>
        <string>do shell script "/bin/bash '/Users/andrew/Documents/Research/AI_workflows/INDD_dashboard/scripts/cron_refresh.sh' &gt;&gt;$HOME/Library/Logs/INDD_dashboard_refresh.out.log 2&gt;&gt;$HOME/Library/Logs/INDD_dashboard_refresh.err.log"</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/andrew/Documents/Research/AI_workflows/INDD_dashboard</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/andrew/Library/Logs/INDD_dashboard_refresh.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/andrew/Library/Logs/INDD_dashboard_refresh.err.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
